# deploy/nifty_iron_fly.service
[Unit]
Description=Nifty Monthly Iron Fly Live Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/nifty-monthly-straddle
ExecStart=/usr/bin/python3 /home/ubuntu/nifty-monthly-straddle/main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nifty_iron_fly
EnvironmentFile=/home/ubuntu/nifty-monthly-straddle/.env
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
